# roles/bridges/templates/systemd.timer
[Unit]
Description={{ item }} bridge timer

[Timer]
OnCalendar={{ roq_config.service.start }}

[Install]
WantedBy=multi-user.target
