# roles/bridges/templates/systemd-stop.timer
[Unit]
Description={{ item }} bridge stop timer

[Timer]
OnCalendar={{ roq_config.service.stop }}

[Install]
WantedBy=multi-user.target
