# roles/roq-json-bridge/templates/roq-json-bridge.timer
[Unit]
Description=roq-{{ config.flags.name }} timer

[Timer]
OnCalendar={{ config.service.start }}

[Install]
WantedBy=multi-user.target
