# roles/roq-json-bridge/templates/roq-json-bridge-stop.timer
[Unit]
Description=roq-{{ config.flags.name }} stop timer

[Timer]
OnCalendar={{ config.service.stop }}

[Install]
WantedBy=multi-user.target
