# roles/conda/templates/systemd-stop.timer
# !!! THIS FILE HAS BEEN AUTO-GENERATED !!!

[Unit]
Description={{ item }} stop timer

[Timer]
OnCalendar={{ config_conda.systemd.timers.stop }}

[Install]
WantedBy=multi-user.target
