# roles/conda/templates/systemd-stop.service
# !!! THIS FILE HAS BEEN AUTO-GENERATED !!!

[Unit]
Description={{ item }} stop service
After=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/echo ''

[Install]
WantedBy=multi-user.target
