# roles/container/templates/systemd.service
# !!! THIS FILE HAS BEEN AUTO-GENERATED !!!

[Unit]
{% if 'systemd' in config_container.services[item] %}
Description={{ config_container.services[item].systemd['description'] | default(item) }}
{% endif %}
Wants=network.target                                                                                             
After=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=2

[Service]
Type={{ config_systemd.profiles[config_container.services[item].systemd.profile].type }}
TimeoutStartSec=0
Restart={{ config_systemd.profiles[config_container.services[item].systemd.profile].restart }}
ExecStop=-/usr/bin/{{ config_container.solution }} stop %n
ExecStartPre=-/usr/bin/{{ config_container.solution }} stop %n
ExecStartPre=-/usr/bin/{{ config_container.solution }} rm %n
ExecStart=/usr/bin/{{ config_container.solution }} run --rm --name %n \
{% if config_container.solution == 'podman' %}
  --cgroup-manager=cgroupfs \
{% endif %}
{% for k,v in config_container.services[item].container.options.items() %}
{% if v is iterable and v is not string and v is not mapping %}
{% for v2 in v %}
  --{{ k }} {{ v2 }} \
{% endfor %}
{% else %}
  --{{ k }} {{ v }} \
{% endif %}
{% endfor %}
  {{ config_container.images[config_container.services[item].container.image].registry }}/{{ config_container.images[config_container.services[item].container.image].name }}:{{ config_container.images[config_container.services[item].container.image].tag }}{{ "" if config_container.services[item].container.args | length == 0 else " \\" }}
{% for k,v in config_container.services[item].container.args.items() %}
{% if v is iterable and v is not string and v is not mapping %}
{% set outer_loop = loop %}
{% for v2 in v %}
  --{{ k }}={{ v2 }}{{ "" if outer_loop.last and loop.last else " \\" }}
{% endfor %}
{% else %}
  --{{ k }}={{ v }}{{ "" if loop.last else " \\" }}
{% endif %}
{% endfor %}

[Install]
WantedBy=multi-user.target
