# roles/roq-license-manager/templates/roq-license-manager.timer
[Unit]
Description=roq-license-manager timer

[Timer]
OnCalendar={{ config.service.start }}

[Install]
WantedBy=multi-user.target
