# roles/services/templates/systemd.timer
# !!! THIS FILE HAS BEEN AUTO-GENERATED !!!

[Unit]
Description={{ item }} timer

[Timer]
OnCalendar={{ conda_config.systemd.timers.start }}

[Install]
WantedBy=multi-user.target
