# roles/services/templates/systemd.timer
# !!! THIS FILE HAS BEEN AUTO-GENERATED !!!

[Unit]
Description={{ item }} timer

[Timer]
OnCalendar={{ roq_config.systemd.start }}

[Install]
WantedBy=multi-user.target
