# roles/services/templates/systemd-stop.timer
# !!! THIS FILE HAS BEEN AUTO-GENERATED !!!

[Unit]
Description={{ item }} stop timer

[Timer]
OnCalendar={{ roq_config.systemd.stop }}

[Install]
WantedBy=multi-user.target
