# roles/gateways/templates/systemd.timer
[Unit]
Description={{ item }} gateway timer

[Timer]
OnCalendar={{ roq_config.service.start }}

[Install]
WantedBy=multi-user.target
